# install/webcdr.service.example
[Service]
Environment=PORT=8080
Environment=CONFIG=cdr
ExecStart=/usr/bin/node /opt/webcdr/serverLdap.js
Restart=always
User=webcdr
Group=webcdr
